# webapp/cron.d/delete_outdated_users
SHELL=/bin/bash
*/5 * * * * root cd /app && /usr/local/bin/pipenv run flask cronjob delete_outdated_users &> /app/cronjob_output
